# scripts/cron.d/pijuice
# every thirty minutes
*/30 * * * * root systemctl restart pijuice.service
